# scripts/stock-analyzer.service
# ==========================================
# stock-analyzer systemd 服务配置
# 安装: sudo cp scripts/stock-analyzer.service /etc/systemd/system/
# 启动: sudo systemctl enable stock-analyzer && sudo systemctl start stock-analyzer
# 查看: sudo systemctl status stock-analyzer
# 日志: journalctl -u stock-analyzer -f
# ==========================================

[Unit]
Description=stock-analyzer A股短线策略分析工具
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/stock-analyzer/backend
Environment="PYTHONPATH=/root/stock-analyzer/backend"
Environment="LLM_API_KEY=sk-82a9eafeccd2482bacdc5b1747f63d7a"
Environment="LLM_BASE_URL=https://api.deepseek.com/v1"
Environment="LLM_MODEL=deepseek-chat"
ExecStart=/root/stock-analyzer/backend/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
